# deploy/roboto.service
# ============================================================
# Roboto — systemd unit (Linux / VPS)
#
# Instalação:
#   sudo cp deploy/roboto.service /etc/systemd/system/roboto.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable roboto
#   sudo systemctl start roboto
#
# Comandos úteis:
#   sudo systemctl status roboto
#   sudo journalctl -u roboto -f        # logs em tempo real
#   sudo systemctl restart roboto
#   sudo systemctl stop roboto
# ============================================================

[Unit]
Description=Roboto — Bot de trading 24/7
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=roboto

# Diretório do projeto (ajuste se necessário)
WorkingDirectory=/opt/roboto

# Carrega variáveis do .env
EnvironmentFile=/opt/roboto/.env

# Comando de inicialização
ExecStart=/opt/roboto/.venv/bin/python -m backend.core.bot

# Restart automático
Restart=always
RestartSec=10

# Redireciona stdout/stderr para journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=roboto

# Limites de segurança
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
